# digital_signage/templates/default/edu.biola.digsig.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>edu.biola.digsig</string>
  <key>KeepAlive</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Digital Signage.app/Contents/MacOS/Digital Signage</string>
    <% unless node['digital_signage']['client']['feed_url'].empty? %>
    <string><%= node['digital_signage']['client']['feed_url'] %></string>
    <% end %>
  </array>
  <key>UserName</key>
  <string>digsig</string>
</dict>
</plist>
