# digital_signage/templates/default/edu.biola.killDigitalSignage.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Enabled</key>
	<true/>
	<key>Label</key>
	<string>edu.biola.killDigitalSignage</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/killall</string>
		<string>Digital Signage</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>4</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
</dict>
</plist>
